# dyndns/debian/postinst
#!/bin/sh
/usr/sbin/update-rc.d atomiadns-dyndns defaults
