# bind_sync/debian/atomiadns-bindsync.service
[Unit]
Description=AtomiaDns Bindsync service

[Service]
Type=simple
ExecStart=/usr/bin/atomiabindsync updated 2>&1 | logger -t atomiabindsync

[Install]
WantedBy=multi-user.target
