# bind_sync/debian/atomiadns-bindsync.service
[Unit]
Description=AtomiaDns Syncer service

[Service]
Type=simple
EnvironmentFile=/etc/atomiadns.conf
ExecStart=/bin/atomiadnssync updated 2>&1 | logger -t atomiadnssync_updated

[Install]
WantedBy=multi-user.target
