# webapp/debian/atomiadns-webapp.postinst
#!/bin/sh

service atomiadns-webapp start
if [ -f /usr/bin/lsb_release ] && [ `/usr/bin/lsb_release -rs | rev | cut -d"." -f2- | rev` -ge "16" ]; then
	#enables systemd service to start on boot up
	systemctl enable atomiadns-webapp
fi